# hades.service
[Unit]
Description=Hades SOC V2.0 - Headless Sentinel Service
Documentation=https://github.com/ManicDevs/hades-soc-v2
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/hades
ExecStart=/opt/hades/sentinel
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hades-sentinel

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/hades
ReadWritePaths=/var/log/hades

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment variables
Environment=HADES_LOG_LEVEL=info
Environment=HADES_DB_PATH=/opt/hades/hades_sentinel.db

[Install]
WantedBy=multi-user.target
